# batch/crontab
* * * * * root /usr/src/app/runbatch.sh >> /var/log/cron.log
